# template/ix86/suse-12.2-JeOS/root/lib/systemd/system/console-setup.service
[Unit]
Description=Show Console

[Service]
ExecStart=-/bin/plymouth --hide-splash

[Install]
WantedBy=runlevel3.target
